# build/watchtower.service
[Unit]
Description=Watchtower Worker
Wants=docker.service

[Service]
ExecStartPre=/bin/sleep 30
ExecStart=/usr/bin/docker exec watchtower_watchtower_1 /usr/local/bin/php /var/www/html/scripts/watchtower.php
Restart=always

[Install]
WantedBy=multi-user.target
